# pi-hole/rootfs/etc/services.d/php-fpm/run
#!/usr/bin/with-contenv bashio
# ==============================================================================
# Community Hass.io Add-ons: Pi-hole
# Runs the PHP-FPM daemon
# ==============================================================================

bashio::log.info "Starting PHP-FPM..."
export PHP_ERROR_LOG=/dev/stderr
exec php-fpm7 --nodaemonize
