# pi-hole/rootfs/etc/services.d/pihole-FTL/run
#!/usr/bin/with-contenv bashio
# ==============================================================================
# Community Hass.io Add-ons: Pi-hole
# Runs the Pi-hole FTL daemon
# ==============================================================================
mkdir -p /var/run/pihole
chown pihole:pihole /var/run/pihole

if bashio::debug; then
    exec s6-setuidgid pihole pihole-FTL debug
else
    exec s6-setuidgid pihole pihole-FTL no-daemon
fi
